# scripts/confluence-sync/com.cortex.confluence-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.cortex.confluence-sync</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-c</string>
    <string>cd "$HOME/tools/A-Team" &amp;&amp; exec "$(which node || echo /usr/local/bin/node)" scripts/confluence-sync/daemon.mjs</string>
  </array>

  <key>WorkingDirectory</key>
  <string>REPLACED_BY_INSTALL_SCRIPT</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/tmp/confluence-sync.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/confluence-sync.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>
</dict>
</plist>
